# vantapi.service
[Unit]
Description=Vant API Service

[Service]
Type=simple
ExecStart=npm run vant-api

Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=production
WorkingDirectory=/home/harry/vant-api

[Install]